# docker/local/supervisor/conf.d/consumer_item_event.ini
[program:consumer_item_event]
process_name=%(program_name)s_%(process_num)02d
command=php artisan consume:item_event
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/consumer_item_event.log
